# photosynthesisv11/photosynthesis-main/sap-with-full-liquid-stake-redemption-workflow/dockernet/dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:myprogram]
command=/usr/local/bin/archwayd

[program:cron]
command=crond -f
